# aacs-update.service
[Unit]
Description=Download and install AACS KEYDB.cfg

[Service]
Type=simple
ExecStart=/opt/aacs-updater.sh

[Install]
Alias=aacs.service
